# recipes/kubernetes-podmaster/podmaster.conf.upstart
# Kubernetes Podmaster (Upstart unit)
description "Kubernetes Podmaster"
start on started kube-apiserver
stop on stopping kube-apiserver

env PODMASTER="/usr/bin/podmaster"
env OPTIONS=""

script
    # read settings like GOMAXPROCS from "/etc/default/consul", if available.
    [ -e /etc/default/podmaster ] && . /etc/default/podmaster

    #export GOMAXPROCS=`nproc`
    exec $PODMASTER $OPTIONS >> /var/log/podmaster.log 2>&1
end script

respawn
respawn limit 10 10
kill timeout 10
